# packages/6-emacs/rdm.service
[Unit]
Description=rdm
Documentation=rdm

[Service]
Type=forking
ExecStart=rdm --daemon
ExecStop=killall rdm
Restart=on-failure

[Install]
WantedBy=default.target

